# 5.3/supervisord.conf
[supervisord]
nodaemon=true

[program:mysql]
command=service mysql start

[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
